# compose/orchestrator/crontab
0 2 * * * cd /opt/tizim && python -m app.main run-daily > /proc/1/fd/1 2>/proc/1/fd/2
